# Docker/Lab 91 Setting Up a Host-Like Environment Using Docker Containers/Docker file Code/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:mysql]
command=/usr/sbin/mysqld
autorestart=true
